# deploy/systemd/fallback-ip.service
[Unit]
Description=Static fallback IP address (192.168.1.250)
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/opt/leidsa-dashboard/deploy/scripts/set-fallback-ip.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
